# data/init/upstart/guerillabackup-generator.conf
# guerillabackup - Start the backup generator service

description "Guerillabackup backup generator service"

start on filesystem
stop on starting rcS

respawn
exec /usr/bin/gb-backup-generator
